# global/etc/xdg/autostart/touchegg.desktop
[Desktop Entry]
Exec=touchegg
Name=Touchégg (gestures daemon)
Type=Application
Version=1.0
