# xfce/etc/xdg/autostart/ksuperkey.desktop
[Desktop Entry]
Exec=ksuperkey
Name=KSuperKey
Type=Application
Version=1.0
Name[en_US]=ksuperkey.desktop
